# services/crop_upload.path
# /etc/systemd/system/crop_upload.timer
## THESE DO NOT WORK. systemd path is not recursive
[Unit]
Description=crop upload
After=network-online.target

[Path]
PathChanged=/mnt/remote

[Install]
WantedBy=timers.target
